# deploy/systemd/usls-gs-mvp-api.service
[Unit]
Description=USLS GS MVP API Service
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=uslsapp
Group=uslsapp
WorkingDirectory=/opt/usls-gs-mvp/server
EnvironmentFile=/etc/usls-gs-mvp/server.env
ExecStart=/usr/bin/node dist/src/index.js
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/usls-gs-mvp/uploads
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
